# conf/pg_sqlog.cron
0 0 * * * postgres psql -AXqtc "SELECT sqlog.cache()" > /dev/null 2>&1
